# templates/nextcloud-cron.service.j2
[Unit]
Description=Nextcloud cron.php job

[Service]
ExecStart=podman exec -u www-data {{ nextcloud_name }} /usr/local/bin/php -f /var/www/html/cron.php
KillMode=process